# crontab
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin/nginx
0	*/8	*	*	*	/usr/local/bin/nginx-rotate-session-ticket-keys
