# crontab
0	*/8	*	*	*	/usr/local/bin/nginx-rotate-session-ticket-keys
